# resources/pihub.timer
[Unit]
Description=Runs the pihub.service every 10 minutes

[Timer]
OnCalendar=*:*
Unit=pihub.service

[Install]
WantedBy=basic.target
